# config/inittab
::respawn:/usr/local/bin/start_syslogd
::respawn:/usr/sbin/postfix start-fg >/dev/null 2>&1
::respawn:/usr/sbin/inetd -f > /dev/null 2>&1
::respawn:/usr/bin/stunnel >/dev/null 2>&1
::respawn:/usr/sbin/saslauthd -a sasldb -c -d >/dev/null 2>&1
::respawn:/usr/sbin/php-fpm84 > /run/php-fpm84.log 2>&1
::respawn:/usr/sbin/nginx -g "daemon off;" >/dev/null 2>&1
::respawn:/usr/local/python/root_runner.py
::respawn:/usr/local/python/doms_runner.py
::respawn:/usr/local/bin/start_wsgi_webui
